# cryptopanner/ops/systemd/cryptopanner-uploader.service
# Uploader: S3 upload with manifest-last ordering + infinite backoff (master spec §9.c, §6.b).
[Unit]
Description=CryptoPanner Uploader
After=network-online.target

[Service]
Type=simple
User=cryptopanner
Group=cryptopanner
Environment=TZ=UTC
Environment=JAVA_OPTS=-Dconfig=/etc/cryptopanner/config.yaml
ExecStart=/opt/cryptopanner/current/uploader/bin/uploader
Restart=always
RestartSec=5
WatchdogSec=30

[Install]
WantedBy=multi-user.target
