# cryptopanner/ops/systemd/cryptopanner-sealer.service
# Sealer: hourly merge + REST backfill + manifest generation (master spec §9, §6.b).
[Unit]
Description=CryptoPanner Sealer
After=network-online.target

[Service]
Type=simple
User=cryptopanner
Group=cryptopanner
Environment=TZ=UTC
Environment=JAVA_OPTS=-Dconfig=/etc/cryptopanner/config.yaml
ExecStart=/opt/cryptopanner/current/sealer/bin/sealer
Restart=always
RestartSec=5
WatchdogSec=30

[Install]
WantedBy=multi-user.target
