# cryptopanner/ops/systemd/cryptopanner-agent.service
# Node Agent: /status,/metrics,/restart,/rotation/trigger HTTP surface (master spec §11.c, §6.b).
# Binds to the VPN interface in production (agent.listen_address); reads heartbeats + drives systemctl.
[Unit]
Description=CryptoPanner Node Agent
After=network-online.target

[Service]
Type=simple
User=cryptopanner
Group=cryptopanner
Environment=TZ=UTC
Environment=JAVA_OPTS=-Dconfig=/etc/cryptopanner/config.yaml
ExecStart=/opt/cryptopanner/current/agent/bin/agent
Restart=always
RestartSec=2
WatchdogSec=30

[Install]
WantedBy=multi-user.target
