# infra/aws/systemd/cryptolake-lifecycle-agent.service
[Unit]
Description=CryptoLake Host Lifecycle Agent
Documentation=https://github.com/your-org/cryptolake
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=cryptolake
Group=docker
Environment=LIFECYCLE_LEDGER_PATH=/data/.cryptolake/lifecycle/events.jsonl
ExecStart=/opt/cryptolake/.venv/bin/python /opt/cryptolake/scripts/host_lifecycle_agent.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptolake-lifecycle

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/data/.cryptolake/lifecycle
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
