# infra/vps/cryptolake.service
[Unit]
Description=CryptoLake Data Pipeline
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/cryptolake
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.vps.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.vps.yml down
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
